# rootfs/etc/s6-overlay/s6-rc.d/avahi/run
#!/command/with-contenv sh

exec /usr/sbin/avahi-daemon --no-rlimits
