# rootfs/etc/s6-overlay/s6-rc.d/smbd/run
#!/command/with-contenv sh

exec smbd -F --debug-stdout --no-process-group
