# rootfs/etc/s6-overlay/s6-rc.d/wsdd2/run
#!/command/with-contenv sh

exec /etc/s6-overlay/scripts/wsdd2
