# rootfs/etc/s6-overlay/s6-rc.d/config/up
with-contenv
/etc/s6-overlay/scripts/config
